# quadlets/llm-routing-valkey-lf.container
[Unit]
Description=Valkey LF — Langfuse v3 BullMQ job queue (POD_NAME_PLACEHOLDER)

[Container]
ContainerName=POD_NAME_PLACEHOLDER-valkey-lf
Image=valkey/valkey:9.1.1
Pod=llm-routing.pod
# play-kube injected <pod-name> -> 127.0.0.1 into /etc/hosts; replicate for apps
# (langfuse 3.222+) that resolve the pod hostname for self/intra-pod calls.
AddHost=POD_NAME_PLACEHOLDER:127.0.0.1
Entrypoint=redis-server
Exec=--bind 127.0.0.1 --port VALKEY_LF_PORT_PLACEHOLDER --requirepass REDIS_AUTH_PLACEHOLDER --maxmemory-policy noeviction --loglevel warning
Volume=DATA_ROOT_PLACEHOLDER/redis-lf-data:/data
HealthCmd=valkey-cli -p VALKEY_LF_PORT_PLACEHOLDER -a REDIS_AUTH_PLACEHOLDER ping
HealthInterval=10s
HealthTimeout=2s
HealthOnFailure=kill
HealthRetries=5
HealthStartPeriod=5s

[Service]
Restart=always
TimeoutStartSec=90
